# src/pcscd_ccid.rules
# udev rules for pcscd and CCID readers

# Gemplus PCMCIA Card
#SUBSYSTEMS=="pcmcia", DRIVER=="serial_cs", ACTION=="add", ATTRS{prod_id1}=="Gemplus", ATTRS{prod_id2}=="SerialPort", ATTRS{prod_id3}=="GemPC Card", RUN+="/usr/sbin/pcscd --hotplug" 

# Kobil mIDentity
BUS=="usb", ATTRS{idVendor}=="0d46", ATTRS{idProduct}=="4081", ACTION=="add", RUN+="/usr/sbin/Kobil_mIDentity_switch"

# generic CCID device
BUS=="usb", ATTRS{bInterfaceClass}=="0b", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"

# non CCID generic (InterfaceClass: 0xFF)
# CherrySmartTerminalST2XXX.txt
BUS=="usb", ATTRS{idVendor}=="046a", ATTRS{idProduct}=="003e", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
# DellSK-3106.txt
BUS=="usb", ATTRS{idVendor}=="413c", ATTRS{idProduct}=="2100", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
# MySmartPad.txt
BUS=="usb", ATTRS{idVendor}=="09be", ATTRS{idProduct}=="0002", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
# SCR331-DI-NTTCom.txt
BUS=="usb", ATTRS{idVendor}=="04e6", ATTRS{idProduct}=="5120", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
# SCR331-DI.txt
BUS=="usb", ATTRS{idVendor}=="04e6", ATTRS{idProduct}=="5111", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
# SDI010.txt
BUS=="usb", ATTRS{idVendor}=="04e6", ATTRS{idProduct}=="5121", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
# SPR532.txt
BUS=="usb", ATTRS{idVendor}=="04e6", ATTRS{idProduct}=="e003", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
# Verisign_secure_storage_token.txt
BUS=="usb", ATTRS{idVendor}=="08e6", ATTRS{idProduct}=="1359", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
# Verisign_secure_token.txt
BUS=="usb", ATTRS{idVendor}=="08e6", ATTRS{idProduct}=="ace0", ACTION=="add", RUN+="/usr/sbin/pcscd --hotplug"
